# openanki.service
[Unit]
Description=OpenAnki Backend Server
After=network.target

[Service]
Type=simple
User=ploi
Group=ploi
WorkingDirectory=/home/ploi/checkst.app/checkstbackend
ExecStart=/home/ploi/checkst.app/checkstbackend/server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
# Environment
Environment=PORT=8080
EnvironmentFile=/home/ploi/checkst.app/checkstbackend/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ploi/checkst.app/checkstbackend/data

[Install]
WantedBy=multi-user.target
